# pinger.conf
#!upstart
description "WiRover Reboot Pinger"

# Start the job
exec python /home/wirover/Arduino-Power-Cycler/udpPing.py

# Start on boot
start on net-device-up IFACE=eth1

